# docker/supervisor/conf.d/laravel-schedule.conf
[program:laravel-schedule]
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
user=root
stdout_logfile=/var/www/html/storage/logs/schedule.log
